# systemd/deskcon-server.service
[Unit]
Description=DeskCon Server Unit

[Service]
Type=simple
ExecStart=INSTALL_DEST_DIR/deskcon.sh

[Install]
WantedBy=multi-user.target
